# systemd/forti-poller.service
[Unit]
Description=FortiGate Monitor Poller
After=network.target

[Service]
Type=simple
User=forti-monitor
Group=forti-monitor
WorkingDirectory=/opt/forti-monitor
EnvironmentFile=-/etc/forti-monitor/env
ExecStart=/opt/forti-monitor/venv/bin/python -m poller.main
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forti-poller

[Install]
WantedBy=multi-user.target
